# gost.service
[Unit]
Description=gost
After=network-online.target
Wants=network-online.target systemd-networkd-wait-online.service

[Service]
Type=simple
User=root
Restart=always
RestartSec=1
DynamicUser=true
LimitNOFILE=200000
ExecStart=/usr/bin/gost -L sni://127.0.0.1:443 -L sni://127.0.0.1:80 -F socks://tw1.dnsunlock.com:10001

[Install]
WantedBy=multi-user.target
